# docker/supervisor.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND
stdout_logfile=/dev/fd/1
stderr_logfile=/dev/fd/2
